# bangle_sleep_stream/receiver/systemd/sleepstream-receiver.service
# File summary:
# Systemd service unit for the Sleep Stream receiver daemon.
# Starts the Python process at boot and restarts on failure.

[Unit]
Description=Sleep Stream BLE Receiver
After=bluetooth.target network.target
Wants=bluetooth.target

[Service]
Type=simple
WorkingDirectory=/opt/bangle_sleep_stream/receiver
ExecStart=/usr/bin/python3 /opt/bangle_sleep_stream/receiver/sleep_receiver.py --name-prefix BangleSleep --db-path /opt/bangle_sleep_stream/receiver/sleepstream.db --log-path /opt/bangle_sleep_stream/receiver/sleepstream.log
Restart=always
RestartSec=3
User=pi
Group=pi

[Install]
WantedBy=multi-user.target
